# service/unitiao-kvm-host.service
[Unit]
Description=UniTao KVM Host
After=network.target libvirtd.service

[Service]
Type=simple
User=root
WorkingDirectory={{REPO_ROOT}}
ExecStart={{REPO_ROOT}}/src/runpy.sh src/rest/app.py --config {{CONFIG_DIR}}/config.json
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
